# waters-core/infrastructure/systemd/integrator-dm.service
[Unit]
Description=WATERS Integrator Data Manager — SCP pipeline 167→238
Documentation=https://github.com/waters-ai/waters-core
After=network.target kafka.service
Wants=kafka.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/WATERS/repos/waters-core
ExecStart=/usr/bin/python3 /home/ubuntu/WATERS/repos/waters-core/agents/integrator_data_manager.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/WATERS/repos/waters-core/logs/integrator_dm.log
StandardError=append:/home/ubuntu/WATERS/repos/waters-core/logs/integrator_dm.log

[Install]
WantedBy=multi-user.target
