# infrastructure/systemd/kapelka-auditor.service
[Unit]
Description=WATERS Kapelka Auditor — tests and reports
Documentation=https://github.com/waters-ai/waters-core
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/WATERS/repos/waters-core/missions/kapelka
ExecStart=/usr/bin/python3 agents/agent_auditor.py
Restart=always
RestartSec=10
StandardOutput=append:logs/auditor.log
StandardError=append:logs/auditor.log
Environment=PYTHONPATH=/home/ubuntu/WATERS/repos/waters-core/missions/kapelka/agents
Environment=KAPELKA_MODE=file
Environment=OLLAMA_HOST=localhost

[Install]
WantedBy=multi-user.target
